# source/examples/services/synx-daemon.service
[Unit]
Description=Synx Code Validation Daemon
Documentation=https://github.com/A5873/synx
After=network.target
Wants=network.target

[Service]
Type=notify
ExecStart=/usr/bin/synx daemon start --config /etc/synx/daemon.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
User=synx
Group=synx

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/synx /var/run/synx
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
